# cm/dot_config/systemd/user/syncthing.service.d/override.conf.tmpl
{{- if ne .chezmoi.hostname "chopin" -}}
{{-   fail (print "Assertion failed: only 'chopin' machine needs the synchting override configuration.") -}}
{{- end -}}
[Unit]
RequiresMountsFor=/mnt/btr
